# server/systemd/dict-english.service
[Unit]
Description=Dict English backend API
After=network.target

[Service]
Type=simple
User=dict-english
Group=dict-english
Environment=DEPLOY_PATH=REPOSITORY_PATH_ON_SERVER
Environment=PYTHONUNBUFFERED=1
ExecStart=sh -lc 'set -a; . "$DEPLOY_PATH/english/.env"; set +a; cd "$DEPLOY_PATH/english/backend" && exec python3 -m vocabdb serve-api --db "$DICT_ENGLISH_DB" --host "$DICT_ENGLISH_HOST" --port "$DICT_ENGLISH_PORT"'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
